# scripts/spmd.service
[Unit]
Description=SPM Daemon — privileged package operations
Documentation=man:spmd(8)
After=network.target local-fs.target

[Service]
Type=simple
ExecStart=/usr/local/bin/spmd
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
ProtectSystem=full
ReadWritePaths=/var/lib/spm /var/cache/spm /etc/spm

[Install]
WantedBy=multi-user.target
